# deploy/rf-controller.service
[Unit]
Description=RF Socket Controller
# Ordering hint only: when a local Mosquitto broker is installed (deploy with
# SETUP_MOSQUITTO=1), start after it so the broker is ready first. This is a
# soft After= with no Requires/Wants — the MQTT client retries and reconnects
# on its own, so a missing or external broker is fine and won't block startup.
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=claw
WorkingDirectory=/home/claw/rf-socket-controller
EnvironmentFile=-/home/claw/rf-socket-controller/.env
ExecStart=/home/claw/rf-socket-controller/rf-controller
Restart=on-failure
RestartSec=5

# Hardening — relax these if you need broader access.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/claw/rf-socket-controller/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
