# deploy/rf-controller.service
[Unit]
Description=RF Socket Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claw
WorkingDirectory=/home/claw/rf-socket-controller
EnvironmentFile=-/home/claw/rf-socket-controller/.env
ExecStart=/home/claw/rf-socket-controller/rf-controller
Restart=on-failure
RestartSec=5

# Hardening — relax these if you need broader access.
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/claw/rf-socket-controller/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
