# deploy/ollama.service.d/override.conf
# systemd drop-in for the Ollama service installed by ollama.com/install.sh.
# Installed to /etc/systemd/system/ollama.service.d/override.conf by
# scripts/setup-ollama.sh, which substitutes @OLLAMA_MODELS@ with the chosen
# model directory (e.g. a 1 TB drive mounted at /mnt/storage).
[Service]
# Bind to loopback only — the assistant is reached via the controller, never
# directly from the LAN.
Environment=OLLAMA_HOST=127.0.0.1:11434
# Keep model blobs off the SD card and on the big drive.
Environment=OLLAMA_MODELS=@OLLAMA_MODELS@
# Keep the model resident for 30 min so back-to-back questions skip the
# brutal cold-reload penalty on a Pi 4.
Environment=OLLAMA_KEEP_ALIVE=30m
# One model in RAM at a time — protects the 8 GB ceiling.
Environment=OLLAMA_MAX_LOADED_MODELS=1
# Don't start until the model drive is mounted.
RequiresMountsFor=@OLLAMA_MODELS@
